# twaine-frontend.service
[Unit]
Description=twAIne Frontend Server
After=network.target twaine-backend.service
Wants=twaine-backend.service

[Service]
Type=simple
User=ale
Group=ale
WorkingDirectory=/home/ale/GIT/twaine
Environment="NODE_ENV=production"
EnvironmentFile=/home/ale/GIT/twaine/.env.local
ExecStartPre=/usr/bin/npm run build
ExecStart=/usr/bin/npm run preview
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=twaine-frontend

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ale/GIT/twaine

[Install]
WantedBy=multi-user.target
