# scripts/twaine-service
#!/bin/bash
### BEGIN INIT INFO
# Provides:          twaine
# Required-Start:    $network $remote_fs $syslog
# Required-Stop:     $network $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: twAIne interactive story platform
# Description:       Starts the twAIne frontend and backend services
### END INIT INFO

# Configuration
APP_NAME="twaine"
APP_DIR="/home/ale/git/twaine"
SERVER_DIR="$APP_DIR/server"
USER="ale"
GROUP="ale"

# PID files
FRONTEND_PID="/var/run/twaine-frontend.pid"
BACKEND_PID="/var/run/twaine-backend.pid"

# Log files
LOG_DIR="/var/log/twaine"
FRONTEND_LOG="$LOG_DIR/frontend.log"
BACKEND_LOG="$LOG_DIR/backend.log"

# Node and npm paths (adjust if needed)
NODE_BIN=$(which node)
NPM_BIN=$(which npm)

# Colors for output
RED='\033[0;31m'
GREEN='\033[0;32m'
YELLOW='\033[1;33m'
NC='\033[0m' # No Color

# Ensure log directory exists
mkdir -p "$LOG_DIR"
chown -R "$USER:$GROUP" "$LOG_DIR"

# Function to check if a process is running
is_running() {
    local pid_file=$1
    if [ -f "$pid_file" ]; then
        local pid=$(cat "$pid_file")
        if ps -p "$pid" > /dev/null 2>&1; then
            return 0
        else
            rm -f "$pid_file"
            return 1
        fi
    fi
    return 1
}

# Function to start the backend server
start_backend() {
    echo -n "Starting twAIne backend server... "
    
    if is_running "$BACKEND_PID"; then
        echo -e "${YELLOW}already running${NC}"
        return 0
    fi
    
    cd "$SERVER_DIR"
    
    # Check if build exists, if not build it
    if [ ! -d "dist" ]; then
        echo "Building backend..."
        su - "$USER" -c "cd $SERVER_DIR && $NPM_BIN run build" >> "$BACKEND_LOG" 2>&1
    fi
    
    # Start the backend server
    su - "$USER" -c "cd $SERVER_DIR && $NODE_BIN dist/index.js >> $BACKEND_LOG 2>&1 & echo \$!" > "$BACKEND_PID"
    
    sleep 2
    
    if is_running "$BACKEND_PID"; then
        echo -e "${GREEN}OK${NC}"
        return 0
    else
        echo -e "${RED}FAILED${NC}"
        return 1
    fi
}

# Function to start the frontend
start_frontend() {
    echo -n "Starting twAIne frontend... "
    
    if is_running "$FRONTEND_PID"; then
        echo -e "${YELLOW}already running${NC}"
        return 0
    fi
    
    cd "$APP_DIR"
    
    # Build and serve the frontend
    su - "$USER" -c "cd $APP_DIR && $NPM_BIN run build >> $FRONTEND_LOG 2>&1"
    su - "$USER" -c "cd $APP_DIR && $NPM_BIN run preview >> $FRONTEND_LOG 2>&1 & echo \$!" > "$FRONTEND_PID"
    
    sleep 2
    
    if is_running "$FRONTEND_PID"; then
        echo -e "${GREEN}OK${NC}"
        return 0
    else
        echo -e "${RED}FAILED${NC}"
        return 1
    fi
}

# Function to stop a service
stop_service() {
    local name=$1
    local pid_file=$2
    
    echo -n "Stopping twAIne $name... "
    
    if ! is_running "$pid_file"; then
        echo -e "${YELLOW}not running${NC}"
        return 0
    fi
    
    local pid=$(cat "$pid_file")
    kill "$pid" 2>/dev/null
    
    # Wait for process to stop
    local count=0
    while ps -p "$pid" > /dev/null 2>&1 && [ $count -lt 10 ]; do
        sleep 1
        count=$((count + 1))
    done
    
    if ps -p "$pid" > /dev/null 2>&1; then
        kill -9 "$pid" 2>/dev/null
        sleep 1
    fi
    
    rm -f "$pid_file"
    echo -e "${GREEN}OK${NC}"
    return 0
}

# Function to show status
show_status() {
    echo "twAIne Service Status:"
    echo "====================="
    
    if is_running "$BACKEND_PID"; then
        local pid=$(cat "$BACKEND_PID")
        echo -e "Backend:  ${GREEN}running${NC} (PID: $pid)"
    else
        echo -e "Backend:  ${RED}stopped${NC}"
    fi
    
    if is_running "$FRONTEND_PID"; then
        local pid=$(cat "$FRONTEND_PID")
        echo -e "Frontend: ${GREEN}running${NC} (PID: $pid)"
    else
        echo -e "Frontend: ${RED}stopped${NC}"
    fi
    
    echo ""
    echo "Log files:"
    echo "  Backend:  $BACKEND_LOG"
    echo "  Frontend: $FRONTEND_LOG"
}

# Main script logic
case "$1" in
    start)
        echo "Starting twAIne services..."
        start_backend
        start_frontend
        echo "Done."
        ;;
    
    stop)
        echo "Stopping twAIne services..."
        stop_service "frontend" "$FRONTEND_PID"
        stop_service "backend" "$BACKEND_PID"
        echo "Done."
        ;;
    
    restart)
        $0 stop
        sleep 2
        $0 start
        ;;
    
    status)
        show_status
        ;;
    
    logs)
        if [ "$2" == "backend" ]; then
            tail -f "$BACKEND_LOG"
        elif [ "$2" == "frontend" ]; then
            tail -f "$FRONTEND_LOG"
        else
            echo "Usage: $0 logs {backend|frontend}"
            exit 1
        fi
        ;;
    
    *)
        echo "Usage: $0 {start|stop|restart|status|logs {backend|frontend}}"
        exit 1
        ;;
esac

exit 0
